# contrib/garuda-low-ram-llm/systemd/ollama.service.d/50-memory-ceiling.conf
# Limit Ollama RAM before the desktop wedges (8 GiB class laptops with iGPU carve-out).
# Tune MemoryMax if you only run very small models (e.g. lower to 3500M).
[Service]
MemoryHigh=3800M
MemoryMax=4800M
