# contrib/garuda-low-ram-llm/systemd/ollama.service.d/60-ollama-concurrency.example.conf
# Optional: align the Ollama *daemon* with tight cgroup RAM (use alongside 50-memory-ceiling.conf).
# Copy as a higher-numbered drop-in so it loads after the memory ceiling, then:
#   sudo systemctl daemon-reload && sudo systemctl restart ollama
#
# Keeps one model resident and one generation at a time inside the service, so MemoryMax
# trips less often while the Study Plan app still serializes its own HTTP calls.

[Service]
Environment=OLLAMA_MAX_LOADED_MODELS=1
Environment=OLLAMA_NUM_PARALLEL=1
